# scrapeNews/install/ansible/roles/scraper/templates/scraper.service.j2
[Unit]
Description=Scraper Service
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory={{ deploy_base }}
ExecStart={{ deploy_base }}/start

[Install]
WantedBy=multi-user.target
